# rc.local
#! /bin/sh -e
//path-to-your-bash-script-<eg>-/home/sammy-ros/.startup/auto_start.sh
exit 0
